# supervisord.conf
[program:chessfellows]
environment: %(host_envs)s
command: gunicorn chessfellows.wsgi:application
directory: /home/ubuntu
autostart: true
stdout_logfile: /tmp/django_log.txt
stderr_redirect: true